# systemd/matrixhub-db-backup.timer
[Unit]
Description=Nightly backup timer for MatrixHub DB

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
